# container-files/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command = /opt/remi/php70/root/sbin/php-fpm --fpm-config /etc/php-fpm.conf
autorestart = true
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log
