# container-files/etc/supervisor.d/www-watch.conf
[program:www-watch]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /data/www; do chown -Rf www:www /data/www; done'
